# services/disk2iso.service
[Unit]
Description=disk2iso - Automatische ISO Erstellung von optischen Medien
After=multi-user.target
Wants=systemd-udevd.service
After=systemd-udevd.service

[Service]
Type=simple
User=root
Group=root
ExecStart=/opt/disk2iso/services/disk2iso/daemon.sh
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

# Nutzt DEFAULT_OUTPUT_DIR aus config.sh
# Ausgabeverzeichnis wurde bei Installation konfiguriert

[Install]
WantedBy=multi-user.target
